# packaging/debian/postinst
#!/bin/sh

set -e

ONE_USER=oneadmin
ONE_GROUP=oneadmin

if [ "$1" = "configure" ]; then
    # Ensure /etc/one directory exists
    mkdir -p /etc/one 2>/dev/null || true
    chgrp "${ONE_GROUP}" /etc/one 2>/dev/null || true
    chmod 0751 /etc/one 2>/dev/null || true

    # Fix permissions for configuration file
    if [ -f /etc/one/cognit-optimizer.conf ]; then
        chown root:"${ONE_GROUP}" /etc/one/cognit-optimizer.conf 2>/dev/null || true
        chmod 0640 /etc/one/cognit-optimizer.conf 2>/dev/null || true
    fi

    # Fix permissions for application directory
    if [ -d /usr/lib/one/cognit-optimizer ]; then
        chown -R root:"${ONE_GROUP}" /usr/lib/one/cognit-optimizer 2>/dev/null || true
        chmod 0750 /usr/lib/one/cognit-optimizer 2>/dev/null || true
        find /usr/lib/one/cognit-optimizer -type f -exec chmod 0640 {} \; 2>/dev/null || true
        find /usr/lib/one/cognit-optimizer -type f -name "*.py" -exec chmod 0750 {} \; 2>/dev/null || true
    fi

    # Fix permissions for venv directory
    if [ -d /usr/share/one/cognit-optimizer ]; then
        chown -R root:"${ONE_GROUP}" /usr/share/one/cognit-optimizer 2>/dev/null || true
        chmod -R 0755 /usr/share/one/cognit-optimizer 2>/dev/null || true
    fi

    # Reload systemd and enable/start the service
    systemctl daemon-reload 2>/dev/null || true
    systemctl enable opennebula-cognit-optimizer 2>/dev/null || true
    systemctl start opennebula-cognit-optimizer 2>/dev/null || true
fi

#DEBHELPER#


